# systemd-units/plymouth-switch-root-initramfs.service.in
[Unit]
Description=Tell Plymouth To Jump To initramfs
DefaultDependencies=no
# dracut-shutdown.service restores the initramfs when it is _stopped_
# use Conflicts to make sure its ExecStop has run before we do
Conflicts=dracut-shutdown.service
After=plymouth-halt.service plymouth-reboot.service plymouth-poweroff.service plymouth-kexec.service dracut-shutdown.service
ConditionPathExists=/run/initramfs/bin/sh

[Service]
Type=oneshot
RemainAfterExit=yes
# Emit a blank message to force a dynamic load of all graphics libraries
# required for rendering text since it's a non-trivially large set, too large
# to manage manually in the initramfs. This will allow casper to render text
# like the live media ejection message after the jump...
ExecStart=-@PLYMOUTH_CLIENT_DIR@/plymouth message --text=""
ExecStart=-@PLYMOUTH_CLIENT_DIR@/plymouth update-root-fs --new-root-dir=/run/initramfs
Type=oneshot
RemainAfterExit=yes
