# systemd-units/systemd-ask-password-plymouth.service.in
[Unit]
Description=Forward Password Requests to Plymouth
Documentation=http://www.freedesktop.org/wiki/Software/systemd/PasswordAgents
DefaultDependencies=no
Conflicts=shutdown.target
After=plymouth-start.service
Before=shutdown.target
ConditionKernelCommandLine=!plymouth.enable=0
ConditionVirtualization=!container
ConditionPathExists=@plymouthruntimedir@/pid

[Service]
ExecStart=@SYSTEMD_ASK_PASSWORD_AGENT@ --watch --plymouth
