# data/indicator-messages.conf.in
description "Indicator Messages Service"

start on indicator-services-start
stop on desktop-end or indicator-services-end

respawn
respawn limit 2 10

exec @pkglibexecdir@/indicator-messages-service

