# data/indicator-datetime.conf.in
description "Indicator Date & Time Backend"

start on indicator-services-start
stop on desktop-end or indicator-services-end

respawn
respawn limit 2 10

exec @pkglibexecdir@/indicator-datetime-service
