# src/libcloudproviders.service.in
[Unit]
Description=Cloud storage providers desktop integration daemon

[Service]
ExecStart=@libexecdir@/cloudprovidersd
Type=dbus
BusName=org.freedesktop.CloudProviders.Daemon
